# .history/BackEnd/supervisor/kanastra-boletos_20240228134223.conf
[program:kanastra-boletos]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --verbose --queue=ImportEstruturaProjectJob --timeout=7200
autostart=true
autorestart=true
startretries=5
startsecs=0
user=root
numprocs=1
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB